# Project/ansible/roles/monitoring/files/node_exporter.service
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=node_exporter
ExecStart=/usr/local/bin/node_exporter --web.listen-address="0.0.0.0:9100"
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
